# checkupdate/checkupdate.service
[Service]

ExecStart=/usr/local/bin/checkupdated.sh

[Install]
WantedBy=multi-user.target
